# tracker/tracker.service
[Unit]
Description=SHIFT Tracker Service

[Service]
User=shift
WorkingDirectory=/home/shift/tracker/bin
ExecStart=/home/shift/tracker/bin/tracker.sh
Restart=always
KillSignal=SIGQUIT
Type=notify
StandardOutput=file:/var/log/tracker.log
StandardError=file:/var/log/tracker.log
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target

